# .config/systemd/user/scream.service
[Unit]
Description=Scream on gateway

[Service]
Type=simple
ExecStart=/sbin/scream -v -i enp6s0

[Install]
WantedBy=default.target
